# docker/runit/qdrouterd/run
#!/bin/bash

exec /sbin/setuser root qdrouterd >> /var/log/qdrouterd.log 2>&1
